# etc/rc.d/rc.local
cd /root && /usr/bin/python2 web.py &>web.log & 
